# systemd/templates/dictation.service.in
# ~/.config/systemd/user/dictation.service
[Unit]
Description=Whisper hotkey dictation
After=whisper.service ydotoold.service
PartOf=speech.target

[Service]
EnvironmentFile=-%h/.config/local-speech/dictation.env
Environment=YDOTOOL_SOCKET=%t/ydotool.sock
ExecStart=__DICTATION_EXEC__ "__LOCAL_SPEECH_ROOT__/dictation.py"
Restart=on-failure
RestartSec=2

[Install]
WantedBy=speech.target
