# systemd/templates/kokoro-tts.service.in
# ~/.config/systemd/user/kokoro-tts.service
[Unit]
Description=Kokoro TTS server (GPU)
After=network.target

[Service]
ExecStart=__KOKORO_EXEC__ "__LOCAL_SPEECH_ROOT__/Kokoro-FastAPI/start-gpu.sh"
WorkingDirectory=__LOCAL_SPEECH_ROOT__/Kokoro-FastAPI
Restart=on-failure
RestartSec=5
Environment=HOME=%h

[Install]
WantedBy=default.target
