# systemd/user/dictation.service
# ~/.config/systemd/user/dictation.service
[Unit]
Description=Whisper hotkey dictation
After=whisper.service ydotoold.service
PartOf=stt.target

[Service]
EnvironmentFile=-%h/.config/local-speech/dictation.env
Environment=YDOTOOL_SOCKET=%t/ydotool.sock
ExecStart=/usr/bin/env mise exec uv -- uv run %h/Projects/local-speech/dictation.py
Restart=on-failure
RestartSec=2

[Install]
WantedBy=stt.target
