# systemd/user/speech.target
[Unit]
Description=Local speech stack
Wants=whisper.service ydotoold.service dictation.service
After=graphical-session.target

[Install]
WantedBy=default.target
