# systemd/user/whisper.service
# ~/.config/systemd/user/whisper.service
[Unit]
Description=whisper.cpp STT server (CUDA)
PartOf=stt.target

[Service]
ExecStart=%h/Projects/local-speech/whisper.cpp/build/bin/whisper-server \
  --model %h/Projects/local-speech/whisper.cpp/models/ggml-small.en.bin \
  --host 127.0.0.1 \
  --port 8080 \
  --convert \
  --inference-path "/v1/audio/transcriptions"
Restart=on-failure

[Install]
WantedBy=stt.target
